# systemd/system/gpuslackbot.service.d/gpuslackbot.conf
[Service]
# TODO: Change these to your SLACK_APP_TOKEN and SLACK_BOT_TOKEN
Environment="SLACK_APP_TOKEN=xapp-***"
Environment="SLACK_BOT_TOKEN=xoxb-***"